# docker/common/supervisor/conf.d/00-portal-backend.conf
[program:portal-backend]
command=/opt/xforce-ai/bin/xforce-portal serve --host 127.0.0.1 --port 8080
autostart=true
autorestart=true
startsecs=2
stdout_logfile=/tmp/xforce-ai/services/portal-backend/stdout.log
stderr_logfile=/tmp/xforce-ai/services/portal-backend/stderr.log
