# docker/common/supervisor/conf.d/10-example-pty-service.conf
[program:example-pty-service]
command=/opt/xforce-ai/bin/xforce-pty-wrap run --run-id example-pty-service --log-dir /tmp/xforce-ai/pty -- python3 -m http.server 18080 --directory /workspace
autostart=false
autorestart=false
stdout_logfile=/tmp/xforce-ai/services/example-pty-service/stdout.log
stderr_logfile=/tmp/xforce-ai/services/example-pty-service/stderr.log
