# docker/common/supervisor/conf.d/20-caddy.conf
[program:caddy]
command=/usr/bin/caddy run --config /tmp/xforce-ai/caddy/Caddyfile.generated --adapter caddyfile
autostart=false
autorestart=true
startsecs=2
stdout_logfile=/tmp/xforce-ai/services/caddy/stdout.log
stderr_logfile=/tmp/xforce-ai/services/caddy/stderr.log
