# docker/common/supervisor/conf.d/30-cloudflared.conf
[program:cloudflared]
command=/opt/xforce-ai/bin/xforce-tunnel auto --url http://127.0.0.1:8088
autostart=false
autorestart=true
startsecs=2
stdout_logfile=/tmp/xforce-ai/services/cloudflared/stdout.log
stderr_logfile=/tmp/xforce-ai/services/cloudflared/stderr.log
